# scripts/legacy/old_configs/systemd/.config/systemd/user/health-check.service
[Unit]
Description=Weekly Health Check
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bin/check-health.sh
